# dictation_remote.service
[Unit]
Description=Whisper Simple Dictation Service

[Service]
WorkingDirectory=%h/apps/whisper-simple-dictation
ExecStart=%h/apps/whisper-simple-dictation/.venv/bin/python dictation.py remote en
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target